# mootools-builder.crontab
MAILTO=oblivious+mtnet@subtlegradient.com
0	*	*	*	*	"$HOME/lib/mootools-builder"
